# install_scripts/etc/rc.local

# CAMERA MODULE
indmod /usr/local/telescope/lib/modules/fliusb.ko

