# install_scripts/profile.d/talon.sh
# Add talon variables
export TELHOME=/usr/local/telescope
export PATH=${TELHOME}/bin:${PATH}
